# pi/deploy/vendiqo-pi-update.service
[Unit]
Description=Update Vendiqo Raspberry Pi edge containers
Requires=docker.service
After=network-online.target docker.service vendiqo-pi.service
Wants=network-online.target

[Service]
Type=oneshot
WorkingDirectory=/opt/vendiqo/pi
EnvironmentFile=-/opt/vendiqo/pi/.env
ExecStart=/opt/vendiqo/pi/pi-ota-update.sh
