# pi/deploy/vendiqo-pi.service
[Unit]
Description=Vendiqo Raspberry Pi edge stack
Requires=docker.service
After=network-online.target docker.service
Wants=network-online.target

[Service]
Type=oneshot
RemainAfterExit=yes
WorkingDirectory=/opt/vendiqo/pi
# Start last-known-good local images only. Delayed OTA (pull/apply) runs via
# vendiqo-pi-update.timer (OnBootSec=5min) using pi-ota-update.sh.
ExecStart=/usr/bin/docker compose -f docker-compose.prod.yml up -d
ExecStop=/usr/bin/docker compose -f docker-compose.prod.yml down
TimeoutStartSec=120

[Install]
WantedBy=multi-user.target
