# pi/deploy/vendiqo-pi-update.service
[Unit]
Description=Update Vendiqo Raspberry Pi edge containers
Requires=docker.service
After=network-online.target docker.service
Wants=network-online.target

[Service]
Type=oneshot
WorkingDirectory=/opt/vendiqo/pi
ExecStart=/usr/bin/docker compose -f docker-compose.prod.yml pull
ExecStart=/usr/bin/docker compose -f docker-compose.prod.yml up -d
